# crontab
0 8,21 * * * root runuser -l pi -c 'cd ~/ilp-commander/ && git pull' && cd /home/pi/ilp-commander/ && sh install_or_update.sh
