# supervisor.conf
[program:ilp-commander]
command=sh -c "python /home/pi/ilp-commander/poller.py || (sleep 10s && false)"
directory=/home/pi/ilp-commander/
startsecs=10
startretries=999999
