# supervisor.conf
[program:ilp-commander]
command=python /home/pi/ilp-commander/poller.py
directory=/home/pi/ilp-commander/
