# supervisor.conf
[program:ilp-commander]
command=/home/pi/.pyenv/shims/python3.6 /home/pi/ilp-commander/poller.py
directory=/home/pi/ilp-commander/
startsecs=20
startretries=999999
